# deploy/tencent-lighthouse/systemd/codewhale-feishu-bridge.service
[Unit]
Description=CodeWhale Feishu/Lark Phone Bridge
Wants=network-online.target codewhale-runtime.service
After=network-online.target codewhale-runtime.service

[Service]
Type=simple
User=codewhale
Group=codewhale
WorkingDirectory=/opt/codewhale/bridge
# Legacy /etc/deepseek is loaded first for old installs; /etc/codewhale wins.
EnvironmentFile=-/etc/deepseek/feishu-bridge.env
EnvironmentFile=-/etc/codewhale/feishu-bridge.env
ExecStart=/usr/bin/node /opt/codewhale/bridge/src/index.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/codewhale-feishu-bridge

[Install]
WantedBy=multi-user.target
